# piksi_multi_cpp/install/base_station.service
[Unit]
Description=Start base station automatically on startup.

[Service]
ExecStart=roslaunch piksi_multi base.launch
